# deploy/systemd/hourei-mcp.service
[Unit]
Description=HOUREI-Assist MCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/shun/HOUREI-Assist
ExecStart=/home/shun/HOUREI-Assist/.venv/bin/hourei-mcp serve
Restart=on-failure
RestartSec=5
Environment=HOUREI_HOST=127.0.0.1
Environment=HOUREI_PORT=8793

[Install]
WantedBy=default.target
